# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-ds5000.postinst
depmod -a
/usr/local/bin/platform_api_mgnt.sh install
systemctl enable pddf-platform-init.service
systemctl start pddf-platform-init.service
systemctl enable syncd-monitor.service
systemctl start syncd-monitor.service
