# files/image_config/script-execution-hooks/post-syncd-start.timer
[Unit]
Description=Timer to run a set of commands after syncd has been initialized
StopPropagatedFrom=syncd.service swss.service

[Timer]
OnActiveSec=3min

[Install]
WantedBy=syncd.service
