# files/image_config/script-execution-hooks/post-syncd-start.service
[Unit]
Description=Timer to run a set of commands after syncd has been initialized
Requisite=syncd.service swss.service

[Service]
Type=oneshot
ExecStart=/bin/true
