# dockers/docker-ptf/conf.d/ssh-host-keys.conf
# Start runtime SSH host-key generation ahead of sshd. If key generation is
# still in progress when sshd first starts, sshd has autorestart enabled and
# will retry after the keys are created.
[program:ssh_host_keys]
command=/usr/local/sbin/docker-ptf-ssh-host-keys.sh
process_name=ssh_host_keys
priority=1
stdout_logfile=/tmp/ssh_host_keys.out.log
stderr_logfile=/tmp/ssh_host_keys.err.log
redirect_stderr=false
autostart=true
autorestart=unexpected
startsecs=0
startretries=3
numprocs=1
