# platform/marvell-prestera/sonic-platform-alliedtelesis/debian/sonic-platform-alliedtelesis-x980-32dq.postinst
#!/bin/sh
# postinst script for sonic-platform-alliedtelesis-x980-32dq
#
# see: dh_installdeb(1)
depmod -a
systemctl enable pddf-platform-init.service
systemctl start pddf-platform-init.service
