# platform/broadcom/sonic-platform-modules-nexthop/nh-4010/etc/systemd/system/opennsl-modules.service.d/ge-override.conf
[Service]
# Disable systemd opennsl-modules.service
# opennsl is initialized through asic_init.sh
ExecStart=
ExecStart=/bin/true
