# files/image_config/cron.d/syslog_counter_update
# Attempts to update syslog counter every 10 minutes
*/10 * * * * root /usr/local/bin/update_syslog_counter.py > /dev/null 2>&1