# platform/marvell-prestera/sonic-platform-asterfusion/debian/sonic-platform-asterfusion-cx532p-n.postinst
#!/bin/sh

python3 -m pip install /usr/share/sonic/device/x86_64-asterfusion_cx532p_n-r0/sonic_platform-1.0-py3-none-any.whl
systemctl enable cx532p-n-platform-init.service
systemctl start cx532p-n-platform-init.service
