# files/image_config/udev/rules.d/99-pmon-dev-bridge.rules
# Tag hot-pluggable bus/FRU device classes with the "pmon-hotplug" udev tag, giving pmon access
# to devices that appear after its container is created. The tag is read by
# get_pmon_device_cgroup_rules() (per-class cgroup grant) and pmon-dev-bridge.sh (mirrors nodes
# into the running pmon container).

# Tag each hot-pluggable bus/FRU class.
SUBSYSTEM=="i2c-dev",                          TAG+="pmon-hotplug"
SUBSYSTEM=="spidev",                           TAG+="pmon-hotplug"
SUBSYSTEM=="uio",                              TAG+="pmon-hotplug"
SUBSYSTEM=="mtd",    KERNEL=="mtd[0-9]*",      TAG+="pmon-hotplug"
SUBSYSTEM=="gpio",   KERNEL=="gpiochip[0-9]*", TAG+="pmon-hotplug"
KERNEL=="ipmi[0-9]*",                          TAG+="pmon-hotplug"
KERNEL=="cpld[0-9]*",                          TAG+="pmon-hotplug"
KERNEL=="fpga[0-9]*",                          TAG+="pmon-hotplug"

# Forward add/remove of a tagged node to the bridge.
TAG=="pmon-hotplug", ACTION=="add|remove", ENV{DEVNAME}!="", RUN+="/usr/local/bin/pmon-dev-bridge.sh event $env{ACTION} $env{DEVNAME}"
