# files/build_templates/sonic.target
[Unit]
Description=SONiC services target.
