# files/build_templates/syslog-counter.service
[Unit]
Description=Syslog counter update Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update_syslog_counter.py