# platform/aspeed/sonic-platform-modules-micas/M2-W6950-128OC/service/micas-usb-network-init.service
[Unit]
Description=Micas USB Network Gadget Service
After=platform_process.service systemd-modules-load.service local-fs.target
Wants=platform_process.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/micas-usb-network-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
