# platform/broadcom/sonic-platform-modules-fs/common/service/tx_temp.service
[Unit]
Description=transceiver temperature Service
After=platform_driver.service subnetwork.service
Requires=platform_driver.service subnetwork.service

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/tx_temp.py

[Install]
WantedBy=multi-user.target
